# Library/LaunchAgents/homebrew.mxcl.spacebar2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>homebrew.mxcl.spacebar2</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/opt/spacebar/bin/spacebar</string>
    <string>--config</string>
    <string>/Users/eduardo/.config/spacebar/spacebarrc2</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>

    <key>USER</key>
    <string>spacebar2</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/spacebar2.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/spacebar2.err.log</string>
</dict>
</plist>
